# deployment/systemd/mempool-api.service
[Unit]
Description=Bitcoin Mempool Fee Prediction API
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/chelo/antigravity/btc/bitcoin-onchain-framework
Environment="PYTHONPATH=/home/chelo/antigravity/btc/bitcoin-onchain-framework"
ExecStart=/home/chelo/antigravity/btc/bitcoin-onchain-framework/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 1234
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
