# deployment/systemd/mempool-collector.service
[Unit]
Description=Bitcoin Mempool Collector Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/chelo/antigravity/btc/bitcoin-onchain-framework
Environment="PYTHONPATH=/home/chelo/antigravity/btc/bitcoin-onchain-framework"
ExecStart=/home/chelo/antigravity/btc/bitcoin-onchain-framework/venv/bin/python /home/chelo/antigravity/btc/bitcoin-onchain-framework/scripts/collector_daemon.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
